# buildroot-external/board/rpi5-kodi-ota/overlay/usr/lib/systemd/system/ota-updater.timer
[Unit]
Description=Periodic OTA update check

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
Persistent=true
Unit=ota-updater.service

[Install]
WantedBy=timers.target
